# contrib/cron/cron.d/etherpad
*/3	*	*	*	*	etherpad /opt/etherpad/contrib/bin/etherpad-watcher.sh
